# vendor/phwoolcon/deploy-automator/example/etc/supervisor/conf.d/deployer.conf
[program:deployer]
# Deploy Automator
directory=/srv/http/deploy.example.com/
command=php bin/cli queue:consume deployer --ttl=600
autorestart=true
user=www-data
environment=SHELL="/bin/bash"
stdout_logfile=/var/log/supervisor/deployer.log
stderr_logfile=/var/log/supervisor/deployer.error.log
